# systemd/ops-heartbeat.timer
[Unit]
Description=OpenClaw Ops Heartbeat Timer (every 5 min)

[Timer]
OnBootSec=60
OnUnitActiveSec=300
AccuracySec=30

[Install]
WantedBy=timers.target
